# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=su postgres -c "postgres -D /var/lib/postgresql/data"
autostart=true
autorestart=true
stdout_logfile=/var/log/postgresql.log
stderr_logfile=/var/log/postgresql-error.log

[program:prometheus]
command=/opt/prometheus/prometheus --config.file=/opt/prometheus/prometheus.yml --storage.tsdb.path=/var/lib/prometheus --web.listen-address=0.0.0.0:9090
autostart=true
autorestart=true
stdout_logfile=/var/log/prometheus.log
stderr_logfile=/var/log/prometheus-error.log

[program:node_exporter]
command=node_exporter --web.listen-address=0.0.0.0:9100
autostart=true
autorestart=true
stdout_logfile=/var/log/node_exporter.log
stderr_logfile=/var/log/node_exporter-error.log

[program:grafana]
command=/opt/grafana/bin/grafana-server --config=/opt/grafana/conf/custom.ini --homepath=/opt/grafana
autostart=true
autorestart=true
stdout_logfile=/var/log/grafana.log
stderr_logfile=/var/log/grafana-error.log

[program:backend]
command=npm start
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend-error.log